# files/crontab
*/5 * * * * python3 /upload/wakeup.py
*/5 * * * * python3 /upload/check.py
